# templates/containers/gotify.container
[Unit]
Description=Gotify notification service
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=gotify
Image=docker.io/gotify/server:2.5.0
Environment=TZ=UTC
PublishPort=127.0.0.1:9090:80
Volume=gotify-data:/app/data
PodmanArgs=--pull=never --read-only=false --pids-limit=256 --memory=256m
Memory=128M
PidsLimit=256

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
