# templates/containers/uptime-kuma.container
[Unit]
Description=Uptime Kuma monitoring
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=uptime-kuma
Image=docker.io/louislam/uptime-kuma:2
PublishPort=10.66.66.1:3001:3001
Volume=/home/podmin/.local/share/uptime-kuma:/app/data:Z
AutoUpdate=registry

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
