# templates/systemd/system/archarden-os-report.service.tmpl
[Unit]
Description=Send daily Arch update and advisory report via ntfy

[Service]
Type=oneshot
ExecStart=__LIB_DIR__/os_update_report.sh
