# templates/systemd/system/archarden-fail2ban-summary.timer.tmpl
[Unit]
Description=Run Archarden Fail2Ban daily summary notification

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
