# templates/systemd/system/archarden-ntfy-healthcheck.timer.tmpl
[Unit]
Description=Run Archarden ntfy endpoint reachability check daily

[Timer]
OnCalendar=*-*-* 09:10:00
Persistent=true

[Install]
WantedBy=timers.target
