# templates/systemd/system/archarden-ntfy-healthcheck.service.tmpl
[Unit]
Description=Archarden ntfy endpoint reachability check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=__LIB_DIR__/ntfy_health_check.sh
