# templates/systemd/system/wg-quick@wg0.service.d/archarden-dnsmasq-refresh.conf.tmpl
# Managed by archarden. Do not edit.
[Service]
ExecStartPost=/usr/bin/systemctl try-restart dnsmasq.service
